# root/etc/s6-overlay/s6-rc.d/svc-gamescope/run
#!/bin/bash
# shellcheck shell=bash
#
# s6 longrun: gamescope driving the real HDMI display, with Steam as its child.
#
# Runs as ROOT (DRM master requires root; the privileged container grants it).
# gamescope's builtin seatd acquires the seat without systemd logind. Steam is
# launched as `abc` via runuser so game state/ownership on the /config PVC is
# preserved.
#
# Configuration via env (defaults shown):
#   STEAM_GAMESCOPE_CONNECTOR=HDMI-A-2   DRM connector to drive
#   STEAM_GAMESCOPE_WIDTH=1600           output width
#   STEAM_GAMESCOPE_HEIGHT=900           output height
#   STEAM_GAMESCOPE_REFRESH=60           output refresh rate

set -u

# s6's with-contenv shebang restores /run/s6/container_environment which would
# clobber our exports. We source the container env explicitly then override.
[ -d /run/s6/container_environment ] && for f in /run/s6/container_environment/*; do
  [ -f "$f" ] && export "$(basename "$f")=$(cat "$f")"
done

# Persist our overrides into the container env too (so children inherit).
mkdir -p /run/s6/container_environment
printf '%s' /config/.XDG > /run/s6/container_environment/XDG_RUNTIME_DIR
printf '%s' 1 > /run/s6/container_environment/WLR_LIBINPUT_NO_DEVICES

CONNECTOR="${STEAM_GAMESCOPE_CONNECTOR:-HDMI-A-2}"
WIDTH="${STEAM_GAMESCOPE_WIDTH:-1600}"
HEIGHT="${STEAM_GAMESCOPE_HEIGHT:-900}"
REFRESH="${STEAM_GAMESCOPE_REFRESH:-60}"
export PATH="/usr/games:/usr/bin:/bin:/usr/sbin:${PATH}"
export WLR_LIBINPUT_NO_DEVICES=1
export HOME=/config
export XDG_RUNTIME_DIR=/config/.XDG
unset WAYLAND_DISPLAY DISPLAY

# gamescope must own /tmp/.X11-unix for its Xwayland, else it errors out.
mkdir -p /tmp/.X11-unix
chown root:root /tmp/.X11-unix
chmod 1777 /tmp/.X11-unix
rm -f /tmp/.X11-unix/* 2>/dev/null

# Clean stale sockets (the /config PVC persists them across restarts).
rm -f "${XDG_RUNTIME_DIR}"/gamescope-* "${XDG_RUNTIME_DIR}"/wayland-* 2>/dev/null || true

# gamescope creates its socket as "gamescope-0", but Steam's pressure-vessel
# sandbox hardcodes WAYLAND_DISPLAY=wayland-0. Launch gamescope in the
# background, wait for gamescope-0, then symlink wayland-0 -> gamescope-0 so
# the sandboxed streaming subprocess can find the compositor.
# Output goes to a PVC log file for diagnosis (s6 logger not used by this base).
mkdir -p /config/logs
/usr/games/gamescope \
  --backend drm \
  -O "${CONNECTOR}" \
  -W "${WIDTH}" -H "${HEIGHT}" -r "${REFRESH}" \
  --steam \
  -- /usr/bin/gamescope-steam-wrapper > /config/logs/gamescope.log 2>&1 &
GSPID=$!

# Wait for gamescope-0 then create the wayland-0 compatibility symlink.
for i in $(seq 1 40); do
  [ -S "${XDG_RUNTIME_DIR}/gamescope-0" ] && break
  sleep 0.25
done
if [ -S "${XDG_RUNTIME_DIR}/gamescope-0" ]; then
  ln -sf gamescope-0 "${XDG_RUNTIME_DIR}/wayland-0"
  echo "[svc-gamescope] symlinked wayland-0 -> gamescope-0"
fi

wait "$GSPID"

